# launchd/com.turkify.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Turkify daemon — oturum açılışında kalıcı düzeltme sürecini başlatır (Faz 4).

  Kurulum:
    1) Aşağıdaki YOL'ları kendi kullanıcı adına göre kontrol et (varsayılan: ekrem).
    2) cp launchd/com.turkify.daemon.plist ~/Library/LaunchAgents/
    3) launchctl load ~/Library/LaunchAgents/com.turkify.daemon.plist
  Durdurma:
    launchctl unload ~/Library/LaunchAgents/com.turkify.daemon.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.turkify.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ekrem/projects/turkify/.venv/bin/python</string>
        <string>-m</string>
        <string>turkify</string>
        <string>serve</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/turkify-daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/turkify-daemon.err.log</string>
</dict>
</plist>
